# deployment/systemd/helionyx-dev.service
[Unit]
Description=Helionyx Personal Substrate (Development)
After=network.target
RequiresMountsFor=/var/lib/helionyx

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/home/john/repos/helio-mono

# Load environment-specific configuration
EnvironmentFile=/home/john/repos/helio-mono/.env.dev
Environment="ENV=dev"

# Run the service
ExecStart=/home/john/repos/helio-mono/.venv/bin/python services/api/runner.py

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helionyx-dev

# Security (optional hardening)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
